# packaging/launchd/com.opencode.guardian.plist.sample
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.opencode.guardian</string>
    <key>ProgramArguments</key>
    <array>
      <string>{{EXEC_PATH}}</string>
      <string>--config</string>
      <string>{{CONFIG_PATH}}</string>
      <string>run-once</string>
      <string>--json</string>
    </array>
    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
